# files/etc/crontabs/root
30 4 * * * sleep 70 && touch /etc/banner && reboot
